# scripte/dyndns/cloudflare/dyndns.cloudflare.service
# Autor: sookie-dev auf GitHub
# Dateien befinden sich unter https://github.homelab.global
#
# chmod +x /data/scripte/dyndns_cloudflare.service
# ln -s /data/scripte/dyndns_cloudflare.service /etc/systemd/system/dyndns_cloudflare.service
# systemctl enable dyndns_cloudflare.service
#
[Unit]
Description=Cloudflare DDNS Updater Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/dyndns_cloudflare.sh update
ExecStop=/usr/local/bin/dyndns_cloudflare.sh stop
Restart=always
RestartSec=10
User=your_user_here  # Falls nötig, ändere zu deinem Benutzer
Environment=PATH=/usr/bin:/usr/local/bin

[Install]
WantedBy=multi-user.target
